# webmin/lsc/debian/webmin-lsc.cron.d
# check for jobs in the queue
*/10 * * * *	root	/usr/share/webmin/lsc/scheduler-daemon.php >/dev/null
